# deploy/celery-beat.service
[Unit]
Description=Trace Celery Beat Scheduler
After=network.target
Requires=celery-worker.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/trace
EnvironmentFile=/opt/trace/.env
ExecStart=/opt/trace/venv/bin/celery -A celery_worker.celery beat --loglevel=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
